# traccar/rootfs/etc/services.d/traccar/run
#!/command/with-contenv bashio
echo "RUN SCRIPT TEST 12345"
sleep 15
exit 1
# ==============================================================================
# Home Assistant Community Add-on: Traccar
# Runs the Traccar daemon
# ==============================================================================
declare -a options

bashio::log.info "Starting Traccar server..."

# Memory settings
options+=(-Xms256m)
options+=(-Xmx512m)

# Prefer IPv4, because of Hassio general issues with IPv6
options+=("-Djava.net.preferIPv4Stack=true")

# Run the JAR
options+=(-jar tracker-server.jar)

# Configuration file
options+=(/config/traccar.xml)

# Run the Traccar daemon

cd /opt/traccar || bashio::exit.nok "Failed to switch to Traccar directory"
ls -la /opt/traccar/ | bashio::log.info
exec java "${options[@]}"

#!/command/with-contenv bashio
declare -a options
declare status

set -x

bashio::log.info "Starting Traccar server..."

options+=(-Xms256m)
options+=(-Xmx512m)
options+=("-Djava.net.preferIPv4Stack=true")
options+=(-jar tracker-server.jar)
options+=(/config/traccar.xml)

cd /opt/traccar || bashio::exit.nok "Failed to switch to Traccar directory"

bashio::log.info "Working dir: $(pwd)"
bashio::log.info "Directory listing before start:"
ls -lah /opt/traccar || true
ls -lah /opt/traccar/logs || true
ls -lah /config || true

java "${options[@]}" > /tmp/traccar-stdout.log 2> /tmp/traccar-stderr.log
status=$?

bashio::log.error "Traccar exited with status ${status}"
bashio::log.error "==== STDOUT ===="
cat /tmp/traccar-stdout.log || true
bashio::log.error "==== STDERR ===="
cat /tmp/traccar-stderr.log || true

bashio::log.error "==== FILES UNDER /opt/traccar ===="
find /opt/traccar -maxdepth 3 -type f | sort || true

for logfile in \
  /opt/traccar/logs/tracker-server.log \
  /opt/traccar/logs/traccar.log \
  /opt/traccar/logs/wrapper.log \
  /opt/traccar/logs/debug.log
do
  if [ -f "${logfile}" ]; then
    bashio::log.error "==== Contents of ${logfile} ===="
    cat "${logfile}" || true
  fi
done

sleep 20
exit "${status}"
